# alarmbox.service
[Unit]
Description=My Python Daemon

[Service]
ExecStart=/usr/bin/python3 /path/to/your/daemon_script.py
WorkingDirectory=/path/to/your/daemon_script_directory
User=your_username
Group=your_groupname
Restart=always

[Install]
WantedBy=multi-user.target
